# debian/cron.d
#
# Regular cron jobs for the debtrack package
#
42 2 * * *	root	/usr/lib/debtrack/debupdatecheck.sh
